# dist/scripts/piscreen-server-controlpanel.service
[Unit]
Description=piscreen-server PHP Control Panel
After=multi-user.target

[Service]
Type=simple
ExecStart=/usr/bin/python /home/pi/piscreen-server/webserver/startcontrolpanel.py
Restart=on-abort

[Install]
WantedBy=multi-user.target
