# dist/scripts/piscreen-server-api.service
[Unit]
Description=piscreen-server PHP API endpoint
After=multi-user.target

[Service]
Type=simple
ExecStart=/usr/bin/bash /home/pi/webserver/startapiserver.sh
Restart=on-abort

[Install]
WantedBy=multi-user.target
