# dist/files/piscreen.service
[Unit]
Description=PiScreen Server
After=multi-user.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /var/www/data/scripts/piscreen.py
Restart=on-abort

[Install]
WantedBy=multi-user.target
